# docker-lyra/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/home/lyra/mail/supervisord.log
pidfile=/tmp/supervisord.pid

[program:scheduler]
command=python3 /home/lyra/scripts/scheduler.py
directory=/home/lyra
autostart=true
autorestart=true
stdout_logfile=/home/lyra/mail/scheduler_stdout.log
stderr_logfile=/home/lyra/mail/scheduler_stderr.log
